# packaging/packetnet.service
[Unit]
Description=Packet.NET node host (packetnet / pdn)
Documentation=https://github.com/packet-net/packet.net
After=network-online.target
Wants=network-online.target

[Service]
# Type=exec for now (the host is a plain WebApplication with no sd_notify
# wiring). A later builder.Host.UseSystemd() upgrades this to Type=notify.
Type=exec
ExecStart=/opt/packetnet/app/packetnet --config /etc/packetnet/packetnet.yaml
WorkingDirectory=/var/lib/packetnet
User=packetnet
Group=packetnet
Restart=on-failure
RestartSec=5

# systemd creates + chowns /var/lib/packetnet (0750, packetnet:packetnet).
StateDirectory=packetnet
# /run/packetnet — where the bundled LOBBY package's daemon binds its Unix socket
# (session.socketPath in its pdn-app.yaml) and the conventional home for other app
# sockets. systemd creates it at start (packetnet:packetnet) and keeps it writable
# under ProtectSystem=strict; pdn-supervised app daemons inherit the sandbox, so
# without this the socket path could never be bound.
RuntimeDirectory=packetnet
# Keep /run/packetnet across a restart. The in-app self-update (Phase 7) stages its request
# here (update.channel + github-update.json) for the privileged packetnet-update.service
# oneshot to read, and the update RESTARTS this very service — without preserve, systemd would
# remove /run/packetnet on stop and the oneshot could lose the spool mid-update.
RuntimeDirectoryPreserve=yes

# The single-file self-contained binary extracts native libs at startup, and
# .NET puts HOME-based state (data-protection keys etc.) under HOME. Both must
# land somewhere writable, since ProtectSystem=strict makes the rest of the FS
# read-only — point them at the StateDirectory. gcServer=0 keeps memory modest
# on small/Pi targets.
Environment=DOTNET_BUNDLE_EXTRACT_BASE_DIR=/var/lib/packetnet
Environment=HOME=/var/lib/packetnet
Environment=DOTNET_gcServer=0

# Bind privileged ports (<1024) as the non-root service user. A pdn-supervised app
# (the BBS) serves IMAP on 993 and SMTP submission on 465/587 — the standard ports
# iPhone Mail's port-less account setup auto-discovers. Ambient caps are inherited by
# the supervised app daemons and survive NoNewPrivileges (systemd grants them at exec,
# not via setuid). Nothing uses it unless an app is configured on a privileged port.
AmbientCapabilities=CAP_NET_BIND_SERVICE

# --- sandboxing (ProtectSystem=strict makes the FS read-only bar the paths below) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# /var/lib/packetnet: the StateDirectory (pdn.db, bundle extract, HOME state).
# /etc/packetnet: the node rewrites its own YAML conffile through the management
# API (config-write, port management, first-run /setup) — an atomic write (a
# .tmp in the same dir, then rename), so the config dir must be writable to the
# service; ProtectSystem=strict otherwise keeps /etc read-only. postinst chowns
# /etc/packetnet to the service user so the unix perms permit it too.
ReadWritePaths=/var/lib/packetnet /etc/packetnet
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
# AF_NETLINK is required by the embedded Tailscale sidecar (sidecar/tsnet, launched
# as a child so it inherits this sandbox): tsnet's userspace network monitor opens a
# netlink socket to watch interface/route changes — without it the sidecar exits at
# startup. (Tailscale's own systemd unit allows AF_NETLINK for the same reason.)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# A serial/USB TNC (serial-kiss / nino-tnc ports) needs device access: add the
# service user to dialout (`usermod -aG dialout packetnet`) and drop in an
# override with DeviceAllow=/dev/ttyACM0 + SupplementaryGroups=dialout. Loopback
# / kiss-tcp / AXUDP nodes need nothing extra.

[Install]
WantedBy=multi-user.target
