# packaging/packetnet-update.service
# packetnet-update.service — the privileged, detached in-app update job for an
# apt-installed pdn node (Phase 7; docs/node-self-update-design.md).
#
# Started ON DEMAND by the node via `systemctl start --no-block
# packetnet-update.service` — a shipped polkit rule (49-packetnet-update.rules)
# authorizes the unprivileged `packetnet` service user to start THIS unit only.
# Deliberately NOT enabled and with no [Install]: it never runs at boot, only when
# the operator clicks Apply. It is its own unit (not a child of packetnet.service)
# so it survives the service restart the upgrade triggers.
[Unit]
Description=Packet.NET node — apply a targeted apt upgrade (in-app update)
Documentation=https://github.com/m0lte/packet.net
# apt needs the network; only meaningful if ever started during boot ordering.
After=network-online.target
Wants=network-online.target
# Don't let an in-app update fire while the box is shutting down.
RefuseManualStop=true

[Service]
Type=oneshot
# Runs as root (the default) — it drives apt, which dpkg + the package's postinst
# (try-restart) use to replace the binary and bounce packetnet.service.
ExecStart=/usr/lib/packetnet/packetnet-apt-update
# A run is a one-off; don't leave the unit "active" after it exits.
RemainAfterExit=no
